# lab_agent/systemd/bonsai-lab-agent.service
[Unit]
Description=Bonsai autonomous coding agent in the untrusted DF lab
After=network-online.target tailscaled.service
Wants=network-online.target tailscaled.service

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/srv/bonsai-agent
EnvironmentFile=/etc/bonsai-agent/agent.env
ExecStart=/opt/bonsai-lab-agent/venv/bin/bonsai-lab-agent
Restart=always
RestartSec=10
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
